# LabController/cron.daily/beaker_expire_osversion
#! /bin/sh
flags=-m
/usr/sbin/tmpwatch "$flags" 10d /var/lib/beaker/osversion_data
