# LabController/cron.hourly/beaker_expire_distros
#!/bin/sh
exec flock -n /var/run/beaker_expire_distros.cron.lock beaker-expire-distros
