# Server/cron.d/beaker
# clean up stale visit rows
52 3 * * * apache /usr/bin/beaker-cleanup-visits
