# Medusa/lab-controller/cron.daily/expire_distros
#!/usr/bin/python

import sys, os
import xmlrpclib

if __name__ == '__main__':
    url = "http://localhost:25152"
    remote = xmlrpclib.ServerProxy(url,allow_none=True)
    token = remote.login("testing","testing")
    
    distros = remote.get_distros()
    for distro in distros:
        if not os.path.exists(distro['kernel']):
            remote.remove_distro(distro['name'], token)
